# deploy/systemd/openclaw-watchdog.service
[Unit]
Description=OpenClaw task watchdog runner
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
Environment=PATH=/root/.nvm/versions/node/v22.22.1/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
WorkingDirectory=/root/.openclaw/workspace
# Default: emit structured watchdog events and auto-resume Codex for RESUME_TASK.
# Cooldown state is persisted under state/watchdog_codex_resume.json.
ExecStart=/usr/bin/env python3 /root/.openclaw/workspace/task_watchdog_runner.py --json --resume-codex --codex-workdir /root/.openclaw/workspace --codex-cooldown-minutes 15
